# deploy/supervisor.conf
; 工作目录需要在项目根目录下
[program:django]
command=python3 backend/manage.py runserver 0.0.0.0:8001

[supervisord]
childlogdir=.

